# lib/systemd/system/startupbeep.service
[Unit]
Description=Plays a beep sound

[Service]
ExecStart=/usr/bin/startupbeep.sh

[Install]
WantedBy=multi-user.target